# roles/jellyfin/templates/jellyfin-container.service.j2
[Unit]
Description=Jellyfin Media Server (Podman container)
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/run/user/{{ jellyfin_uid }}/containers

[Service]
User=jellyfin
Group=jellyfin
Environment=XDG_RUNTIME_DIR=/run/user/{{ jellyfin_uid }}
Environment=HOME=/home/jellyfin
Type=oneshot
RemainAfterExit=yes
Restart=on-failure
TimeoutStopSec=70

ExecStartPre=-/usr/bin/podman stop jellyfin
ExecStartPre=-/usr/bin/podman rm jellyfin

ExecStart=/usr/bin/podman run \
    --detach \
    --name jellyfin \
    --userns=keep-id \
    --security-opt label=type:jellyfin_container_t \
    --security-opt no-new-privileges \
    --read-only \
    --tmpfs /tmp:rw,nosuid,nodev,size=512m \
    --tmpfs /run:rw,nosuid,nodev,size=64m \
    --cap-drop ALL \
    --memory {{ jellyfin_memory_limit }} \
    -v {{ jellyfin_config_dir }}:/config \
    -v {{ jellyfin_cache_dir }}:/cache \
    -v /srv/media:/media:ro \
    -p {{ jellyfin_http_port }}:8096 \
    -p {{ jellyfin_https_port }}:8920 \
    {{ jellyfin_container_image }}@{{ jellyfin_container_digest }}

ExecStop=-/usr/bin/podman stop jellyfin
ExecStopPost=-/usr/bin/podman rm jellyfin

SystemCallArchitectures=native
RemoveIPC=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
